# showcase.conf
description "showcase"
author  "Dyn Inc"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

exec /usr/local/bin/node /opt/showcase/bin/server.js 2>&1 | logger -p daemon.info -t showcase

